# deploy/systemd/diplom-healthcheck.service
[Unit]
Description=Diplom annotation server healthcheck and repair
After=network.target nginx.service diplom-gunicorn.service postgresql.service

[Service]
Type=oneshot
ExecStart=/home/oldskull/apps/Diplomnaya-rabota/deploy/support/diplom-repair.sh
