# deploy/systemd/diplom-vk-monitoring-test.timer
[Unit]
Description=Run Diplom VK monitoring test collector hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=180

[Install]
WantedBy=timers.target
